# mbsync.service
[Unit]
Description=Sync mailbox

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -V primary
